# deploy/mr-rediscover.service
[Unit]
Description=market-research weekly source rediscovery
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=mr
Group=mr
EnvironmentFile=/etc/mr/env
ExecStart=/usr/local/bin/mr rediscover --all
Nice=10
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/mr
NoNewPrivileges=true
PrivateTmp=true
